# systemd/miflora-read.timer
[Unit]
Description=Run Xiaomi Mi Plant update hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
